# assets/etc/cron.d/haproxy.cron
11 4 * * * root /usr/local/sbin/certbot-renew >/dev/null
11 5 * * * root /usr/local/sbin/certbot-ocsp  >/dev/null
